# files/instavpn.conf
start on runlevel [2345]
stop on runlevel [016]
respawn
script
  cd /opt/instavpn/
  exec node server.js
end script
#manual